# lib/rc.local
#!/bin/sh -e
#
# rc.local
#
# This is the rc.local file that will run everytime the SmartMirror boots and will trigger the initialization script
# 

# Print the IP address
_IP=$(hostname -I) || true
if [ "$_IP" ]; then
  printf "My IP address is %s\n" "$_IP"
fi

cd /home/pi/SmartMirror/lib/scr

./PIinit.sh &

exit 0
